# claude-code/config/com.minitom.cliproxyapi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.minitom.cliproxyapi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/minitom/workspaces/harness-configs/CLIProxyAPI/cli-proxy-api</string>
        <string>--config</string>
        <string>config.yaml</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/minitom/workspaces/harness-configs/CLIProxyAPI</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/tmp/cliproxyapi.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cliproxyapi.stderr.log</string>
</dict>
</plist>
